# systemd/llama-server-duo.service
[Unit]
Description=llama.cpp server (Qwen2.5-Coder-7B for Claude Code, CPU)
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/distrobox enter llama-cpu -- \
    /var/home/%u/llama.cpp/build/bin/llama-server \
    -m /var/home/%u/models/qwen2.5-coder-7b/Qwen2.5-Coder-7B-Instruct-Q4_K_M.gguf \
    --host 127.0.0.1 \
    --port 8080 \
    -c 32768 \
    --parallel 1 \
    -t 16 \
    --jinja \
    --temp 0.7 \
    --top-p 0.8 \
    --top-k 20
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
