# systemd/llama-server-duo.service
[Unit]
Description=llama.cpp server (DeepSeek-Coder-V2-Lite-16B for Claude Code, CPU prebuilt)
After=network.target

[Service]
Type=simple
ExecStartPre=-/usr/bin/podman rm -f llama-server
ExecStart=/usr/bin/podman run --rm --name llama-server -p 127.0.0.1:8080:8080 -v /var/home/%u/models:/models:z ghcr.io/ggml-org/llama.cpp:server -m /models/deepseek-coder-v2-lite/DeepSeek-Coder-V2-Lite-Instruct-Q4_K_M.gguf --host 0.0.0.0 --port 8080 --alias deepseek-coder-v2-lite --parallel 1 -c 65536 -t 16 --jinja
ExecStop=/usr/bin/podman stop -t 5 llama-server
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
